# web/deploy/whipround-web.service
# Whipround web — systemd unit (standalone Astro/Node server)
#
# Install (see DEPLOY.md for the full procedure):
#   sudo cp whipround-web.service /etc/systemd/system/whipround-web.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now whipround-web
#
# Runs `node ./dist/server/entry.mjs` (the @astrojs/node standalone build).
# Secrets (STRIPE_SECRET_KEY, WHIPROUND_STRIPE_PRODUCT_ID) live in the
# EnvironmentFile below, never in this unit and never in git.

[Unit]
Description=Whipround campaign site (Astro/Node standalone)
After=network.target

[Service]
Type=simple
# Where you copied the build output (contains ./dist/server/entry.mjs):
WorkingDirectory=/var/www/whipround
Environment=HOST=127.0.0.1
Environment=PORT=4321
Environment=NODE_ENV=production
EnvironmentFile=-/etc/whipround/web.env
# Use the absolute path from `which node` on the box (nvm paths differ):
ExecStart=/usr/bin/node ./dist/server/entry.mjs
Restart=always
RestartSec=3
User=www-data
Group=www-data
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/whipround
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
